# crontab
45 3 * * 6 /home/alex/Project1/script.sh
