# crontab
echo
45 3 * * 6 crontab >/dev/null 2>&1
